# templates/default/web/uwsgi.upstart.erb
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN

description "uWSGI"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec /usr/local/bin/uwsgi --ini /etc/uwsgi/graphite.ini

#post-start script
#  sleep 1 # giving the daemon a chance to create the socket file
#  chmod -f 777 <%= node.graphite.web.socket %>
#end script

post-stop script
  rm -f <%= node.graphite.web.pidfile %>
  rm -f <%= node.graphite.web.socket %>
end script
